# contrib/heartbeat/esper-heartbeat.timer
[Unit]
Description=ESPer heartbeat

[Timer]
OnCalendar=minutely

[Install]
WantedBy=multi-user.target
